# files/etc/uci-defaults/99-custom
uci -q batch << EOF
set network.lan.ipaddr='192.168.32.253'
EOF
